# SOURCES/check4updates.cron
#! /bin/sh

# Randomize within 5 min
sleep $(expr $RANDOM \% 300)
check4updates -m
